# autostart/compiz.desktop
[Desktop Entry]
Type=Application
Name=compiz
Comment=%COMMENT%
Exec=compiz --replace
Terminal=false
