# autostart/xfce4-config.desktop
[Desktop Entry]
Name=xfce4-config
Exec=bash /usr/share/xfce4/xfce4-config.sh
Type=Application
Icon=
Comment=%COMMENT%
